# templates/default/hastebin.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:             hastebin
# Default-Start:        2 3 4 5
# Default-Stop:
# Short-Description:    Hastebin server
### END INIT INFO

## init script for hastebin nodejs server

dir="/data/hastebin"
user=<%= node[:hastebin][:user] %>
node="/usr/local/bin/node"
log="/var/log/hastebin.log"

start () {
  echo "Starting hastebin ..."
  su - $user -c "cd ${dir} && ${node} ${dir}/server.js > ${log} 2>&1 &"
}

stop () {
  echo "Stopping hastebin ..."
  pkill -f $dir
}

case $1 in

  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 2
    start
    ;;
  *)
    echo "Usage: $0 {stop|start}"
    ;;

esac
